# bin/zenbook-duo-keyboard.service
[Unit]
Description=ASUS Zenbook Duo Keyboard Backlight Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/keyboard-service -c /etc/keyboard-service/hypr_monitor.conf
Restart=on-failure
StandardOutput=journal
StandardError=journal
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
